# Daemon/sentinel.service
# Please install this by moving the .service file to /etc/systemd/system/ 
# Also please move the Daemon script to /usr/bin/
[Unit]
Description=The Sentinel Daemon by Ty Chermsirivatana

[Service]
ExecStart=/usr/bin/mydaemon
Restart=on-failure

[Install]
WantedBy=multi-user.target 